# root/etc/uci-defaults/30_luci-theme-pinky
uci batch <<-EOF
	set luci.themes.Pinky=/luci-static/pinky
	set luci.main.mediaurlbase=/luci-static/pinky
	commit luci
EOF
exit 0
